# lib/hb-browser-editor/hb-browser-editor-idle.timer
[Unit]
Description=Check for idle browser editor periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
